# rootfs/etc/services.d/owserver/run
#!/command/with-contenv bashio

if ! pgrep -x "/opt/owfs/bin/owserver" > /dev/null
then
    bashio::log.info "Starting owserver..."
    exec /opt/owfs/bin/owserver -c /etc/owfs.conf
fi