# rootfs/etc/s6-overlay/s6-rc.d/owhttpd/run
#!/command/with-contenv bashio

if bashio::config.true 'owhttpd' ;then
    if ! pgrep -x "/opt/owfs/bin/owhttpd" > /dev/null
    then
        bashio::log.info "Wait until owserver is up and running..."
        bashio::net.wait_for 4304
        bashio::log.info "Starting owhttpd..."
        exec /opt/owfs/bin/owhttpd -r -c /etc/owfs.conf
    fi
fi